# udev-script/rules.d/mount.rules
# /etc/udev/rules.d/mount.rules
ENV{DEVTYPE}=="disk", ACTION=="add", ENV{ID_BUS}=="usb", RUN+="/usr/local/bin/dumpusb %E{DEVNAME} &"
